# arablife-bot.service
[Unit]
Description=ArabLife Discord Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/ArabLife
ExecStart=/root/ArabLife/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/arablife-bot.log
StandardError=append:/var/log/arablife-bot.error.log

[Install]
WantedBy=multi-user.target
